# roles/webapp/templates/container-webapp-proxy.socket.j2
[Socket]
ListenStream={{ app_proxy_port }}

[Install]
WantedBy=sockets.target
